# debian/postinst
#!/bin/bash
set -e

# Function to check if a VyOS firewall group exists
check_firewall_group_exists() {
    local group_type="$1"
    local group_name="$2"

    # Use show command to check if group exists
    if /opt/vyatta/sbin/vyatta-cfg-cmd-wrapper show firewall group "$group_type" "$group_name" >/dev/null 2>&1; then
        return 0  # Group exists
    else
        return 1  # Group does not exist
    fi
}

case "$1" in
    configure)
        echo "Configuring VyOS IP Blocklist Generator..."

        # Check if we're running on VyOS
        if [ -f "/opt/vyatta/sbin/vyatta-cfg-cmd-wrapper" ]; then
            echo "VyOS detected, configuring firewall groups..."

            # Track if we made any changes
            changes_made=false

            # Create VyOS firewall groups automatically
            /opt/vyatta/sbin/vyatta-cfg-cmd-wrapper begin

            # Check and create IPv4 network group if it doesn't exist
            if check_firewall_group_exists "network-group" "threats-blocklist-ipv4"; then
                echo "IPv4 firewall group 'threats-blocklist-ipv4' already exists, skipping creation"
            else
                echo "Creating IPv4 firewall group 'threats-blocklist-ipv4'"
                /opt/vyatta/sbin/vyatta-cfg-cmd-wrapper set firewall group network-group threats-blocklist-ipv4 description 'IPv4 Threat Intelligence (auto-created)'
                changes_made=true
            fi

            # Check and create IPv6 network group if it doesn't exist
            if check_firewall_group_exists "ipv6-network-group" "threats-blocklist-ipv6"; then
                echo "IPv6 firewall group 'threats-blocklist-ipv6' already exists, skipping creation"
            else
                echo "Creating IPv6 firewall group 'threats-blocklist-ipv6'"
                /opt/vyatta/sbin/vyatta-cfg-cmd-wrapper set firewall group ipv6-network-group threats-blocklist-ipv6 description 'IPv6 Threat Intelligence (auto-created)'
                changes_made=true
            fi

            # Commit the changes only if we made any
            if [ "$changes_made" = true ]; then
                echo "Committing firewall group changes..."
                if /opt/vyatta/sbin/vyatta-cfg-cmd-wrapper commit; then
                    echo "✓ VyOS firewall groups created successfully"
                else
                    echo "⚠ Warning: Failed to commit VyOS configuration"
                    /opt/vyatta/sbin/vyatta-cfg-cmd-wrapper end
                    exit 1
                fi
            else
                echo "No changes needed for firewall groups"
            fi

            /opt/vyatta/sbin/vyatta-cfg-cmd-wrapper end

            echo ""
            echo "VyOS firewall groups created! Next steps:"
            echo "1. Configure firewall rules to use these groups:"
            echo "   set firewall ipv4 forward filter rule 12 source group network-group 'threats-blocklist-ipv4'"
            echo "   set firewall ipv6 forward filter rule 16 source group ipv6-network-group 'threats-blocklist-ipv6'"
            echo "2. Run: vyos-ipblock --verbose"
            echo "3. Run: /config/scripts/sync-vyos-threats.sh"

        else
            echo "Non-VyOS system detected, skipping firewall group creation"
        fi

        # Create directories for configuration files
        mkdir -p /config/scripts
        mkdir -p /var/log/vyos-ipblock

        # Set proper permissions
        chmod 755 /config/scripts
        chmod 755 /var/log/vyos-ipblock

        echo ""
        echo "VyOS IP Blocklist Generator installed successfully!"
        echo "Configuration files: /config/scripts/"
        echo "Logs: /var/log/vyos-ipblock/"
        echo ""
        echo "Quick start:"
        echo "1. Test: vyos-ipblock --dry-run --verbose"
        echo "2. Run: vyos-ipblock --verbose"
        echo "3. Enable automation: systemctl enable --now vyos-ipblock.timer"
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

#DEBHELPER#

exit 0
