# debian/postinst
#!/bin/bash
set -e

case "$1" in
    configure)
        # Create directories
        mkdir -p /config/scripts
        mkdir -p /var/log/vyos-ipblock
        
        # Set permissions
        chmod 755 /config/scripts
        chmod 755 /var/log/vyos-ipblock
        
        # Copy example whitelist if it doesn't exist
        if [ ! -f /config/scripts/whitelist.txt ]; then
            cp /usr/share/doc/vyos-ipblock/examples/whitelist.txt.example /config/scripts/whitelist.txt
            chmod 644 /config/scripts/whitelist.txt
        fi
        
        # Enable systemd timer (but don't start it automatically)
        systemctl daemon-reload
        systemctl enable vyos-ipblock.timer
        
        echo "VyOS IP Blocklist Generator has been installed."
        echo ""
        echo "Configuration:"
        echo "  - Edit /config/scripts/whitelist.txt to add your networks"
        echo "  - Add AbuseIPDB API key to /config/scripts/abuseipdb.key (optional)"
        echo ""
        echo "Usage:"
        echo "  - Test: vyos-ipblock --dry-run --verbose"
        echo "  - Run once: vyos-ipblock"
        echo "  - Enable automatic updates: systemctl start vyos-ipblock.timer"
        echo ""
        echo "Documentation: /usr/share/doc/vyos-ipblock/"
        ;;
esac

#DEBHELPER#

exit 0
